# deploy/supervisor/cbt-worker.conf
; Bare-metal / VM queue workers for the CBT backend.
; Install to /etc/supervisor/conf.d/cbt-worker.conf, then:
;   supervisorctl reread && supervisorctl update && supervisorctl start cbt-worker:*
;
; These workers grade submitted exams off the request thread. They drain the
; `grading` queue first (the end-of-exam auto-submit spike), then `default`
; (notifications etc.). Scale throughput with numprocs.
[program:cbt-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/cbt/backend/artisan queue:work redis --queue=grading,default --sleep=1 --tries=5 --max-time=3600 --max-jobs=1000 --backoff=10
directory=/var/www/cbt/backend
user=www-data
numprocs=4
autostart=true
autorestart=true
; Let an in-flight grade finish on restart rather than be killed mid-transaction.
stopwaitsecs=70
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/cbt/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

; Laravel scheduler — runs exam:submit-expired every minute (the server-side
; exam timer: force-submits sessions whose crashed/tampered client never did).
; Start with: supervisorctl start cbt-scheduler
[program:cbt-scheduler]
command=php /var/www/cbt/backend/artisan schedule:work
directory=/var/www/cbt/backend
user=www-data
autostart=true
autorestart=true
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/cbt/scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
