# backend/docker/supervisord.conf
; One container, three roles: PHP-FPM, nginx, and a pool of queue workers that
; grade submitted exams off the request thread. Run with:
;   supervisord -c /etc/supervisor/conf.d/cbt.conf
[supervisord]
nodaemon=true
user=root
pidfile=/run/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:php-fpm]
command=php-fpm --nodaemonize --force-stderr
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Grading workers. They drain the `grading` queue first (the end-of-exam spike),
; then `default` (notifications etc.). numprocs scales grading throughput; raise
; it for very large cohorts. stopwaitsecs lets an in-flight grade finish on
; shutdown rather than being killed mid-transaction.
[program:queue-worker]
command=php /app/artisan queue:work redis --queue=grading,default --sleep=1 --tries=5 --max-time=3600 --max-jobs=1000 --backoff=10
process_name=%(program_name)s_%(process_num)02d
numprocs=4
autostart=true
autorestart=true
priority=20
stopwaitsecs=70
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
